# root/etc/uci-defaults/luci-canto-migration
#!/bin/sh

ucode "/etc/canto/scripts/migrate_config.uc"

exit 0
