# root/etc/uci-defaults/luci-unison-migration
#!/bin/sh

ucode "/etc/unison/scripts/migrate_config.uc"

exit 0
