# root/etc/init.d/homeproxy
#!/bin/sh /etc/rc.common
# SPDX-License-Identifier: GPL-3.0-only
#
# Copyright (C) 2022 ImmortalWrt.org

USE_PROCD=1

START=99
STOP=10

CONF="homeproxy"
PROG="/usr/bin/sing-box"
TUNNAME="emortal-singbox"

HP_DIR="/etc/homeproxy"
RUN_DIR="/var/run/homeproxy"
LOG_PATH="$RUN_DIR/homeproxy.log"
DNSMASQ_DIR="/tmp/dnsmasq.d"

FW4="$(command -v fw4)"
IPT_RULE_NAME="homeproxy_input_rule"

log() {
	echo -e "$(date "+%Y-%m-%d %H:%M:%S") [DAEMON] $*" >> "$LOG_PATH"
}

setup_server_firewall() {
	local cfg="$1"

	local enabled
	config_get_bool enabled "$cfg" "enabled"
	[ "$enabled" = "1" ] || return 1

	local port
	config_get port "$cfg" "port"

	# TODO: support nftables based firewall4
	if [ -z "$FW4" ]; then
		for i in "iptables" "ip6tables"; do
			[ -n "$(command -v "$i")" ] || continue
			for v in "tcp" "udp"; do
				echo -e "$i -t filter -A $IPT_RULE_NAME -p $v --dport $port -j ACCEPT" >> "$RUN_DIR/$CONF.include"
			done
		done
	fi
}

start_service() {
	config_load "$CONF"

	local routing_mode
	config_get routing_mode "config" "routing_mode" "bypass_mainland_china"

	local main_node="nil"
	local default_outbound="nil"
	if [ "$routing_mode" != "custom" ]; then
		config_get main_node "config" "main_node" "nil"
	else
		config_get default_outbound "routing" "default_outbound" "nil"
	fi

	local server_enabled
	config_get_bool server_enabled "server" "enabled" "0"

	if [ "$main_node" = "nil" ] && [ "$default_outbound" = "nil" ] && [ "$server_enabled" = "0" ]; then
		return 1
	fi

	mkdir -p "$RUN_DIR"
	"$HP_DIR"/scripts/generate_sing-box.lua 2>>"$LOG_PATH"

	local inbounds="$(jsonfilter -i "$RUN_DIR/sing-box.json" -e "@.inbounds" 2>"/dev/null")"
	if [ ! -e "$RUN_DIR/sing-box.json" ]; then
		log "Error: failed to generate configuration."
		exit 1
	elif [ -z "$inbounds" ] || [ "$inbounds" = "[ ]" ]; then
		log "Error: no valid inbound found."
		exit 1
	elif ! "$PROG" check --config "$RUN_DIR/sing-box.json" 2>>"$LOG_PATH"; then
		log "Error: wrong configuration detected."
		exit 1
	fi

	if [ "$main_node" != "nil" ] || [ "$default_outbound" != "nil" ]; then
		# Auto update
		local auto_update auto_update_time
		config_get_bool auto_update "subscription" "auto_update" "0"
		if [ "auto_update" = "1" ]; then
			config_get auto_update_time "subscription" "auto_update_time" "2"
			echo "0 ${auto_update_time} * * * $HP_DIR/scripts/update_geodata.sh update_version update_subscription" >> "/etc/crontabs/root"
		fi
		# Clean log
		echo -e "*/3 * * * * $HP_DIR/scripts/clean_log.sh" >> "/etc/crontabs/root"
		/etc/init.d/cron restart

		# Hijack DNS
		mkdir -p "$DNSMASQ_DIR"
		cat <<-EOF >> "$DNSMASQ_DIR/dnsmasq-homeproxy.conf"
			no-poll
			no-resolv
			server=8.8.8.8
		EOF
		/etc/init.d/dnsmasq restart
	fi

	# Set firewall
	# TODO: support nftables based firewall4
	if [ "$main_node" != "nil" ] || [ "$default_outbound" != "nil" ]; then
		# Allow FORWARD
		if [ -z "$FW4" ]; then
			for i in "iptables" "ip6tables"; do
				[ -n "$(command -v "$i")" ] || continue
				echo -e "$i -I FORWARD -m comment --comment 'Forward HomeProxy TUN' -o '$TUNNAME' -j ACCEPT" >> "$RUN_DIR/$CONF.include"
			done
		fi
	fi
	local server_auto_firewall
	config_get_bool server_auto_firewall "server" "auto_firewall" "0"
	if [ "$server_auto_firewall" = "1" ]; then
		if [ -z "$FW4" ]; then
			for i in "iptables" "ip6tables"; do
				[ -n "$(command -v "$i")" ] || continue
				echo -e "$i -N $IPT_RULE_NAME" >> "$RUN_DIR/$CONF.include"
				echo -e "$i -t filter -I INPUT -j $IPT_RULE_NAME" >> "$RUN_DIR/$CONF.include"
			done
		fi
		config_foreach setup_server_firewall "server"
	fi
	if [ -z "$FW4" ]; then
		sh "$RUN_DIR/$CONF.include"
	fi

	procd_open_instance "sing-box"
	procd_set_param command "$PROG"
	procd_append_param command run --config "$RUN_DIR/sing-box.json"
	procd_set_param stderr 1
	procd_set_param respawn
	procd_close_instance

	echo > "$RUN_DIR/sing-box.log"
	log "$(sing-box version | awk 'NR==1{print $1,$3}') started."
}

stop_service() {
	sed -i "/$CONF/d" "/etc/crontabs/root"
	/etc/init.d/cron restart

	# Revert DNS hijack
	rm -f "$DNSMASQ_DIR/dnsmasq-homeproxy.conf"
	# For backward compatibility only
	[ "$(uci -q get "dhcp.@dnsmasq[0]._unused_homeproxy_changed")" != "1" ] || {
		uci -q del "dhcp.@dnsmasq[0].noresolv"
		uci -q del_list "dhcp.@dnsmasq[0].server"="8.8.8.8"
		uci -q rename "dhcp.@dnsmasq[0]._orig_noresolv"="noresolv"
		uci -q rename "dhcp.@dnsmasq[0]._orig_server"="server"
		uci -q del "dhcp.@dnsmasq[0]._unused_homeproxy_changed"
		uci -q commit "dhcp"
	}
	/etc/init.d/dnsmasq restart

	# TODO: support nftables based firewall4
	if [ -z "$FW4" ] && [ -s "$RUN_DIR/$CONF.include" ]; then
		for i in "iptables" "ip6tables"; do
			[ -n "$(command -v "$i")" ] || continue
			"$i" -D FORWARD -m comment --comment "Forward HomeProxy TUN" -o "$TUNNAME" -j ACCEPT
			"$i" -F "$IPT_RULE_NAME"
			"$i" -t filter -D INPUT -j "$IPT_RULE_NAME"
			"$i" -X "$IPT_RULE_NAME"
		done 2>"/dev/null"

		echo > "$RUN_DIR/$CONF.include"
	fi

	rm -f "$RUN_DIR/sing-box.json" "$RUN_DIR/sing-box.log"

	log "Service stopped."
}

reload_service() {
	log "Reloading service..."

	stop
	start
}

service_triggers() {
	procd_add_reload_trigger "$CONF"
}
