# root/etc/init.d/homeproxy
#!/bin/sh /etc/rc.common
# SPDX-License-Identifier: GPL-3.0-only
#
# Copyright (C) 2022 ImmortalWrt.org

USE_PROCD=1

START=99
STOP=10

CONF="homeproxy"
PROG="/usr/bin/sing-box"
TUNNAME="emortal-singbox"

HP_DIR="/etc/homeproxy"
RUN_DIR="/var/run/homeproxy"
LOG_PATH="$RUN_DIR/homeproxy.log"

FW4="$(command -v fw4)"
IPT_RULE_NAME="homeproxy_inbound_rule"

log() {
	echo -e "$(date "+%Y-%m-%d %H:%M:%S") [DAEMON] $*" >> "$LOG_PATH"
}

setup_server_firewall() {
	local cfg="$1"

	local enabled
	config_get_bool enabled "$cfg" "enabled"
	[ "$enabled" = "1" ] || return 1

	local port
	config_get port "$cfg" "port"

	if [ -z "$FW4" ]; then
		for i in "iptables" "ip6tables"; do
			[ -n "$(command -v "$i")" ] || continue
			for v in "tcp" "udp"; do
				echo -e "$i -t filter -A $IPT_RULE_NAME -p $v --dport $port -j ACCEPT" >> "/var/etc/$CONF.include"
			done
		done
	fi
}

start_service() {
	config_load "$CONF"

	local routing_mode
	config_get routing_mode "config" "routing_mode" "bypass_mainland_china"

	local main_server="nil"
	local default_outbound="nil"
	if [ "$routing_mode" != "custom" ]; then
		config_get default_outbound "routing" "default_outbound" "nil"
	else
		config_get main_server "config" "main_server" "nil"
	fi

	local server_enabled
	config_get_bool server_enabled "server" "enabled" "0"

	if [ "$main_server" = "nil" ] && [ "$default_outbound" = "nil" ] && [ "$server_enabled" = "0" ]; then
		exit 1
	fi

	mkdir -p "$RUN_DIR"
	"$HP_DIR"/scripts/generate_sing-box.lua

	local inbounds="$(jsonfilter -i "$RUN_DIR/sing-box.json" -e "@.inbounds")"
	if [ ! -e "$RUN_DIR/sing-box.json" ]; then
		log "Error: failed to generate configuration."
		exit 1
	elif [ -z "$inbounds" ] || [ "$inbounds" = "[ ]" ]; then
		log "Error: no valid inbound found."
		exit 1
	elif ! "$PROG" check --config "$RUN_DIR/sing-box.json" 2>>"$LOG_PATH"; then
		log "Error: wrong configuration detected."
		exit 1
	fi

	# Auto update
	local auto_update auto_update_time
	config_get_bool auto_update "subscription" "auto_update" "0"
	if [ "auto_update" = "1" ]; then
		config_get auto_update_time "subscription" "auto_update_time" "2"
		echo "0 ${auto_update_time} * * * $HP_DIR/scripts/update_geodata.sh update_version update_subscription" >> "/etc/crontabs/root"
	fi
	# Clean log
	echo -e "*/10 * * * * $HP_DIR/scripts/clean_log.sh" >> "/etc/crontabs/root"
	/etc/init.d/cron restart

	# Hijack DNS
	uci rename "dhcp.@dnsmasq[0].nonwildcard"="_orig_nonwildcard" 2>"/dev/null"
	uci rename "dhcp.@dnsmasq[0].dns_redirect"="_orig_dns_redirect" 2>"/dev/null"
	uci set "dhcp.@dnsmasq[0].nonwildcard"=1
	uci commit "dhcp"

	# Set firewall
	if [ -z "$FW4" ]; then
		mkdir -p "/var/etc/"

		# Allow FORWARD
		for i in "iptables" "ip6tables"; do
			[ -n "$(command -v "$i")" ] || continue
			echo -e "$i -I FORWARD -m comment --comment \"Forward HomeProxy TUN\" -o $TUNNAME -j ACCEPT" >> "/var/etc/$CONF.include"
		done

		local server_auto_firewall
		config_get_bool server_auto_firewall "server" "auto_firewall" "0"
		if [ "$server_auto_firewall" = "1" ]; then
			for i in "iptables" "ip6tables"; do
				[ -n "$(command -v $i)" ] || continue
				echo -e "$i -N $IPT_RULE_NAME" >> "/var/etc/$CONF.include"
				echo -e "$i -t filter -I INPUT -j $IPT_RULE_NAME" >> "/var/etc/$CONF.include"
			done
			config_foreach setup_server_firewall "server"
		fi

		sh "/var/etc/$CONF.include"
	fi

	procd_open_instance "sing-box"
	procd_set_param command "$PROG"
	procd_append_param command run --config "$RUN_DIR/sing-box.json"
	procd_set_param stderr 1
	procd_set_param respawn
	procd_close_instance

	echo > "$RUN_DIR/sing-box.log"
	log "$(sing-box version | awk '{print $1,$2}') started."
}

stop_service() {
	sed -i "/$CONF/d" "/etc/crontabs/root"
	/etc/init.d/cron restart

	uci del "dhcp.@dnsmasq[0].nonwildcard" 2>"/dev/null"
	uci rename "dhcp.@dnsmasq[0]._orig_nonwildcard"="nonwildcard" 2>"/dev/null"
	uci rename "dhcp.@dnsmasq[0]._orig_dns_redirect"="dns_redirect" 2>"/dev/null"
	uci commit "dhcp"

	if [ -z "$FW4" ]; then
		for i in "iptables" "ip6tables"; do
			[ -n "$(command -v $i)" ] || continue
			"$i" -D FORWARD -m comment --comment "Forward HomeProxy TUN" -o "$TUNNAME" -j ACCEPT
			"$i" -F "$IPT_RULE_NAME"
			"$i" -t filter -D INPUT -j "$IPT_RULE_NAME"
			"$i" -X "$IPT_RULE_NAME"
		done 2>"/dev/null"

		rm -f "/var/etc/$CONF.include"
	fi

	rm -f "$RUN_DIR/sing-box.json" "$RUN_DIR/sing-box.log"

	log "Service stopped."
}

reload_service() {
	log "Reloading service..."

	stop
	start
}

service_triggers() {
	procd_add_reload_trigger "$CONF"
}
